# HsTools/UserSystemd/Hs-PartitionMonitor.service
[Unit]
Description=Check for unused mounts

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo /usr/local/bin/Hsr-PartitionMonitor
