# HsTools/UserSystemd/Hs-StatusWarnings.service
[Unit]
Description=Monitor the status of sensitive devices

[Service]
Type=oneshot
ExecStart=/usr/local/bin/Hs-StatusWarnings
